# src/configs/crontab
0 12 * * * /snap/bin/certbot renew --quiet --nginx